# manifest/etc/services.d/mpd/run
#!/usr/bin/execlineb -P
mpd --stdout --no-daemon /etc/mpd.conf
