# manifest/etc/services.d/ympd/run
#!/usr/bin/execlineb -P
ympd --host localhost --port 6600 --webport 8866 --user ympd
